# systemd/wallpaper/wallpaper.timer
[Unit]
Description=Shuffle through wallpapers

[Timer]
OnStartupSec=1
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target